# book/systemd/clang_systemd/helloworld_c.service
[Unit]
Description=HelloWorld Service

[Service]
Type=notify
ExecStart=/usr/local/bin/helloworld_c

[Install]
WantedBy=multi-user.target
