# book/systemd/clang_systemd/helloworld.service
[Unit]
Description=HelloWorld Service

[Service]
Type=notify
ExecStart=/usr/local/bin/helloworld

[Install]
WantedBy=multi-user.target
