# book/erlang_systemd/helloworld.service
[Unit]
Description=HelloWorld Service

[Service]
Type=notify
ExecStart=/home/ericksun/install/otp-26.1/bin/erl -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/rebar/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/cf/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/cth_readable/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/certifi/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/eunit_formatters/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/getopt/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/providers/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/bbmustache/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/ssl_verify_fun/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/relx/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/bootstrap/lib/erlware_commons/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/default/lib/providers/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/default/lib/getopt/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/default/lib/cf/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/default/lib/erlware_commons/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/default/lib/certifi/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/default/lib/rebar/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/cf/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/cth_readable/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/certifi/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/eunit_formatters/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/getopt/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/providers/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/bbmustache/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/ssl_verify_fun/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/relx/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/erlware_commons/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/.erlang.mk/rebar3/_build/prod/lib/rebar/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/deps/systemd/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/deps/enough/ebin -pa /home/ericksun/program/rabbitmq-server-debug/book/erlang_systemd/ebin -noinput  -s erlang_systemd_app boot

[Install]
WantedBy=multi-user.target
